# pikonek.cron
@daily root find /var/log/server.log -mtime +1 -delete >/dev/null 2>&1